# crontab
*/15 8-23 * * * neil /usr/bin/flock -n /bin/bash /home/neil/tdp/tdp-sync.sh